# templates/default/ttyS1.conf.erb
# Generated by Chef.
# Local modifications will be overwritten.

# <%= node['sol']['tty']['name'] %> - getty
#
# This service maintains a getty on <%= node['sol']['tty']['name'] %> from the point the system is
# started until it is shut down again.

start on stopped rc RUNLEVEL=[2345] and (
            not-container or
            container CONTAINER=lxc or
            container CONTAINER=lxc-libvirt)

stop on runlevel [!2345]

respawn
exec /sbin/getty -8 <%= node['sol']['serial']['speed'] %> <%= node['sol']['tty']['name']  %>
